# debian/rclocal.postinst
#! /bin/sh
set -e

echo "postinst called with: $1" >&2

case "$1" in
	configure)
		update-rc.d -f rclocal remove
		update-rc.d rclocal defaults 99 00
		mkdir -p /home/local/adm/rc-start
		mkdir -p /home/local/adm/rc-stop
		;;

	abort-upgrade|abort-remove|abort-deconfigure)
    	;;

	*)
		echo "postinst called with unknown argument \`$1'" >&2
		exit 1
		;;
esac

#DEBHELPER#
exit 0
